# tdm/rootfs/etc/services.d/tdm/finish
#!/usr/bin/execlineb -S1
# shellcheck shell=bash
# ==============================================================================
# Restart tdm if it fails
# ==============================================================================
if [[ "$1" -ne 0 ]] && [[ "$1" -ne 256 ]]; then
  bashio::log.warning "Restarting TDM..."
  /run/s6/bin/s6-svc -t /run/s6/services/tdm
fi