# tdm/rootfs/etc/services.d/tdm/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start TDM
# ==============================================================================
TMP_FOLDER=$(mktemp -d)
VNC_PORT="$(bashio::addon.port 5900)"
ARCH="$(bashio::info.arch)"
DEBUG=$(bashio::config "debug")
export HOME='/share/tdm'
mkdir -p ${HOME}/.vnc
(
    echo "vncconfig --nowin &"
    echo "exec i3"
) > ${HOME}/.vnc/xstartup
chmod 755 ${HOME}/.vnc/xstartup

(
        echo "SecurityTypes=None"
        echo "geometry=1920x1080"
) > ${HOME}/.vnc/config
chmod 755 ${HOME}/.vnc/config

# check if i3 config exists if not copy default
if [ ! -f ${HOME}/.config/i3/config ]; then
    mkdir -p ${HOME}/.config/i3
    cp /tmp/i3-config ${HOME}/.config/i3/config
fi
# Get MQTT config
if bashio::services.available "mqtt"; then
    MQTT_HOST=$(bashio::services "mqtt" "host")
    MQTT_PASSWORD=$(bashio::services "mqtt" "password")
    MQTT_PORT=$(bashio::services "mqtt" "port")
    MQTT_USER=$(bashio::services "mqtt" "username")
    (
        echo "[General]"
        echo "connect_on_startup=true"
        echo "hostname=${MQTT_HOST}"
        echo "password=${MQTT_PASSWORD}"
        echo "port=${MQTT_PORT}"
        echo "username=${MQTT_USER}"
    ) > ${HOME}/TDM/tdm.cfg
fi
# Run it only on localhost if not expose
if bashio::var.has_value "${VNC_PORT}"; then
    bashio::log.warning "Your direct VNC access is not protected!"
    LOCAL_ONLY=no
else
    LOCAL_ONLY=yes
    echo "localhost" >> ${HOME}/.vnc/config
fi

export BROWSER=/usr/bin/firefox-esr
export XDG_RUNTIME_DIR="${TMP_FOLDER}"
export DISPLAY=":0"

bashio::log.info "Initializing VNC server (local/${LOCAL_ONLY})..."
# if debug is enabled, start vnc in debug mode
if [ "${DEBUG}" == "true" ]; then
    bashio::log.info "Starting VNC server in debug mode..."
    /usr/bin/vncserver "$DISPLAY" &
else
    bashio::log.info "Starting VNC server in normal mode..."
    /usr/bin/vncserver "$DISPLAY" >/dev/null 2>&1 &
fi
# Wait for VNC server to start before continuing
bashio::log.info "TDM waiting for VNC to start"
bashio::net.wait_for 5900

# Start TDM
bashio::log.info "Starting tdmgr..."
exec /usr/bin/python3 /usr/bin/tdmgr.py