# tdm/rootfs/etc/services.d/nginx/finish
#!/usr/bin/execlineb -S1
# shellcheck shell=bash
# ==============================================================================
# Take down the S6 supervision tree when nginx fails
# ==============================================================================
if [[ "$1" -ne 0 ]] && [[ "$1" -ne 256 ]]; then
  bashio::log.warning "Halting add-on"
  /run/s6/basedir/bin/halt
fi