# tdm/rootfs/etc/services.d/nginx/run
#!/usr/bin/env bashio
# shellcheck shell=bash
# ==============================================================================
# Start NGINX service
# ==============================================================================

# Wait for NGINX to start before continuing
bashio::net.wait_for 5900

bashio::log.info "Starting Nginx..."
exec nginx
